# src/local-provider/storage-provider-local.service.in
[Unit]
Description=Storage Framework local provider

[Service]
Type=dbus
BusName=com.canonical.StorageFramework.Provider.Local
ExecStart=@CMAKE_INSTALL_FULL_LIBDIR@/@PROJECT_NAME@/storage-provider-local
